# felucca.service
[Unit]
Description = Felucca Service

[Service]
User = jenkins
Type = forking
ExecStart = /var/lib/Felucca/felucca.sh start
ExecStop = /var/lib/Felucca/felucca.sh stop
ReStart =  /var/lib/Felucca/felucca.sh stop && sleep 1 && /var/lib/Felucca/felucca.sh start

[Install]
WantedBy=multi-user.target